# infra/unix/systemd/eventra-cache-loader.service
[Unit]
Description=Eventra Cache Loader
After=network-online.target redis.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/eventra/current
EnvironmentFile=/etc/eventra/eventra.env
ExecStart=/usr/bin/npm --workspace @eventra/cache-loader run start
Restart=always
RestartSec=5
User=root
Group=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
